# net-im/j2j/files/j2j-initd
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

PIDFILE="/var/run/jabber/j2j.pid"

depend() {
    need net
    use jabber-server
}

start() {
    ebegin "Starting J2J Jabber Transport"
    start-stop-daemon \
        --start \
        --pidfile ${P} \
        --chuid jabber:jabber \
        --background --make-pidfile \
        --exec /usr/bin/python INSPATH/main.py \
        -- \
        -c /etc/jabber/j2j.conf
    eend $?
}

stop() {
    ebegin "Stopping J2J Jabber Transport"
    start-stop-daemon --stop --quiet --pidfile ${P}
    eend $?
}
